# daily.crontab
0 7 * * * cd /app && bundle exec ruby service.rb 2>&1 > service.log
